# profile.d/rover.conf.sh
# Rover Configuration 

export SSH_LOGIN_TARGET="https://cupcake.zyp.io/api/sshlogger"
